# modules/oraclejdk/files/profile.d/java.sh
PATH=/usr/java/default/bin:$PATH
JAVA_HOME=/usr/java/default
export PATH JAVA_HOME

