# modules/oraclejdk/files/profile.d/java.csh
set path = ( /usr/java/default/bin $path )
setenv JAVA_HOME /usr/java/default

